# site-cookbooks/aic13/templates/default/follower_fetcher.conf.erb
description "FollowerFetcher"

start on startup
stop on shutdown

respawn

script
  echo $$ > /var/run/follower_fetcher.pid
  cd /home/deploy/follower_fetcher/
  exec sudo -u deploy /home/deploy/follower_fetcher/run_follower_fetcher_0.sh >> /var/log/aic13_follower_fetcher.log
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/aic13_follower_fetcher.log
end script

pre-stop script
    rm /var/run/follower_fetcher.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/aic13_follower_fetcher.log
end script
