# site-cookbooks/aic13/templates/default/webapp.conf.erb
# /etc/init/webapp.conf

description "AIC13-G4-T4 Web App"

start on runlevel [2345]
stop on runlevel [^2345]

# Restart when job dies
respawn

# Give up restart after 5 respawns in 60 seconds
respawn limit 5 60

script

  chdir /home/deploy/weabpp
  exec /usr/bin/node server.js 2>&1

end script